# cftravel-api.service
[Unit]
Description=CFTravel IA Agent API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/iagent.cftravel.net
Environment=ENVIRONMENT=production
Environment=API_HOST=0.0.0.0
Environment=API_PORT=8000
ExecStart=/usr/bin/python3 cftravel_py/api/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 