# cftravel-api.service
[Unit]
Description=CFTravel Python API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/cftravel/cftravel_py
ExecStart=/usr/bin/python -m api.server
Restart=always
RestartSec=10
Environment=PYTHONPATH=/path/to/your/cftravel3/cftravel_py

[Install]
WantedBy=multi-user.target 